# roles/supervisor_setup/files/flaskwpage.conf
[program:flaskwpage]
command=python /var/www/flaskwpage/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/flaskwpage/long.err.log
stdout_logfile=/var/log/flaskwpage/long.out.log